# deploy/sm-adviser-intraday.timer
# Fires every 15 min during market hours (Mon-Fri 09:00-15:45 IST); the job itself gates
# precisely to 09:15-15:30 and no-ops otherwise. Set the NUC tz to Asia/Kolkata.
[Unit]
Description=Schedule SM-Adviser intraday price refresh

[Timer]
OnCalendar=Mon-Fri 09..15:00/15
Persistent=false

[Install]
WantedBy=timers.target
